# app/layout/Library/LaunchDaemons/ch.xxtou.tfcontainermanagerd.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>_AdditionalProperties</key>
	<dict>
		<key>RunningBoard</key>
		<dict>
			<key>Managed</key>
			<false/>
			<key>Reported</key>
			<false/>
		</dict>
	</dict>
	<key>HighPriorityIO</key>
	<true/>
	<key>KeepAlive</key>
	<true/>
	<key>EnvironmentVariables</key>
	<dict>
		<key>PATH</key>
		<string>/usr/local/xxtouch/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin</string>
	</dict>
	<key>StandardOutPath</key>
	<string>/usr/local/xxtouch/log/ch.xxtou.tfcontainermanagerd.stdout.log</string>
	<key>StandardErrorPath</key>
	<string>/usr/local/xxtouch/log/ch.xxtou.tfcontainermanagerd.stderr.log</string>
	<key>WorkingDirectory</key>
	<string>/usr/local/xxtouch</string>
	<key>UserName</key>
	<string>root</string>
	<key>Label</key>
	<string>ch.xxtou.tfcontainermanagerd</string>
	<key>EnableTransactions</key>
	<false/>
	<key>EnablePressuredExit</key>
	<false/>
	<key>ProgramArguments</key>
	<array>
		<string>/usr/local/xxtouch/bin/tfcontainermanagerd</string>
	</array>
	<key>RunAtLoad</key>
	<true/>
	<key>ThrottleInterval</key>
	<integer>5</integer>
	<key>POSIXSpawnType</key>
	<string>Adaptive</string>
	<key>ProcessType</key>
	<string>Adaptive</string>
</dict>
</plist>
